# classiwhale/supervisord.conf
[program:celeryd]
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py celeryd -l info

# [program:rabbitmq-server]
# command=rabbitmq-server
